# rpmbuild/SOURCES/k8s-worker-finalize.service
[Unit]
Description=K8s worker finalize (install containerd/kubelet after rpm txn)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
Environment=K8S_RELEASE=v1.35
ExecStart=/usr/libexec/k8s-worker/finalize.sh

